# crontab
# m h  dom mon dow   command
*/10 * * * * ~/tweet_crawler.py